# launchd/nodev.comments.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>build.nodev.comments</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/jinsik/Desktop/Workspace/threads/scripts/respond.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/jinsik/Desktop/Workspace/threads</string>
    <key>StartInterval</key>
    <integer>1800</integer>
    <key>StandardOutPath</key>
    <string>/Users/jinsik/Desktop/Workspace/threads/logs/comments.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/jinsik/Desktop/Workspace/threads/logs/comments.err</string>
</dict>
</plist>
